# 104_wireguard/image/wireguard_dashboard
#!/sbin/openrc-run
description="Wireguard Dashboard"

depend() {
    need net
    need localmount
    need wireguard0
}

start() {
    cd /opt/WGDashboard/src
    ./wgd.sh start
}

stop() {
    cd /opt/WGDashboard/src
    ./wgd.sh stop
}
